# deployment/MediCureFlow.service
[Unit]
Description=MediCureFlow Healthcare Management System
Documentation=https://github.com/your-org/MediCureFlow
After=network.target postgresql.service redis.service
Wants=postgresql.service redis.service

[Service]
Type=notify
User=www-data
Group=www-data
RuntimeDirectory=MediCureFlow
WorkingDirectory=/opt/MediCureFlow
Environment="PATH=/opt/MediCureFlow/venv/bin"
Environment="DJANGO_SETTINGS_MODULE=MediCureFlow.settings.production"
ExecStart=/opt/MediCureFlow/venv/bin/gunicorn --config gunicorn.conf.py MediCureFlow.wsgi:application
ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
TimeoutStopSec=5
PrivateTmp=true
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=MediCureFlow

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/MediCureFlow/logs /opt/MediCureFlow/media /tmp
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictSUIDSGID=true
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
